# deploy/systemd/predict-settlement.timer
[Unit]
Description=Run Predict settlement cycle every minute

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
AccuracySec=5s
Unit=predict-settlement.service

[Install]
WantedBy=timers.target
